# scripts/emanate-api.service.in
[Unit]
Description=Emanate API daemon
After=network.target

[Service]
Type=simple
WorkingDirectory={{PATH_TO_EMANATE_REST_API}}
ExecStart=/usr/bin/node {{PATH_TO_EMANATE_REST_API}}/index.js
Restart=always

[Install]
WantedBy=multi-user.target

